# SOURCES/cs50-appliance-2.1-5/etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

# snapshots (to avoid premature mounting by systemd)
# http://forums.fedoraforum.org/showpost.php?p=1492420&postcount=5
/bin/mount -o ro,vers=3 -t nfs localhost:/.private/snapshots /snapshots
